# ansible/roles/viking/files/viking-daphne.service
[Unit]
Description=Daphne web server for Viking.
After=network.target

[Service]
Type=forking
User=viking
Group=viking
WorkingDirectory=/home/viking/project/
ExecStart=/home/viking/scripts/viking daphne
Restart=always

[Install]
WantedBy=multi-user.target
